# deploy/tunnel/canopy-tunnel.service
[Unit]
Description=Canopy relay cloudflared tunnel (durable public wss:// for the local relay)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=%h/.local/bin/canopy-tunnel.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
